# etc/crontab
@reboot systemctl start reflector.timer
@hourly systemctl start reflector.timer
